# luci-theme-vbeauty/root/uci-defaults/64_luci-theme-rosy
#!/bin/sh

if [ "$PKG_UPGRADE" != 1 ]; then
	uci batch <<-EOF
		set luci.themes.Vbeauty=/luci-static/vbeauty
		set luci.main.mediaurlbase=/luci-static/vbeauty
		commit luci
	EOF
fi

exit 0